# netperf/config/systemd/tasks/netperf-test-isp.service
[Unit]
Description=Test ISP connection speed

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/netperf/test_network.py isp
