# netperf/config/systemd/tasks/netperf-test-local.service
[Unit]
Description=Test local network connection speeds

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/netperf/test_network.py local

