# netperf/config/systemd/tasks/netperf-test-dns.service
[Unit]
Description=Test DNS queries

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/netperf/test_network.py dns

