# netperf/config/systemd/netperf-interfaces.service
[Unit]
Description = Network Performance Monitor Interface Configuration
Requires = network.target
After = netperf-storage.service network.target

[Service]
Type=forking
ExecStart = /usr/bin/python3 /opt/netperf/configure_interfaces.py
ExecStartPost = /usr/bin/bash -c 'systemctl restart nginx.service'
Type=oneshot

[Install]
WantedBy = multi-user.target
